# systemd/mnt-bpshare-bpshare0-backup.automount.template
[Unit]
Description=Automount for Borg Backup Share

[Automount]
# Must match the Where= in the .mount file
Where=/mnt/YOUR_MOUNTPOINT
# How long to wait before unmounting when idle (in seconds)
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
